# archlinux/sealkey@.path
[Unit]
Description = Watch initrd for changes to reseal TPM key

[Path]
PathChanged = /boot/initramfs-%i.img

[Install]
WantedBy=sysinit.target

